# supervisor.conf
; ─────────────────────────────────────────────────────────────────────────────
; Supervisor config for cuan-workflow Laravel queue workers
; Production Linux / Ubuntu VPS
;
; Install Supervisor:  apt install supervisor
; Enable:              supervisorctl reread && supervisorctl update
; Start all:           supervisorctl start cuan-workflow-worker:*
; Status:              supervisorctl status
; ─────────────────────────────────────────────────────────────────────────────

[program:cuan-workflow-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/cuan-workflow/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --timeout=120 --queue=default
; ↑ Change path to actual deployment path on your VPS

directory=/var/www/cuan-workflow
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data

; ── Scaling: change this number to add more workers ─────────────────────────
; 3 workers = good for 500 concurrent users
; 5 workers = good for 1000+ concurrent users
; Each worker consumes ~30–50MB RAM
numprocs=3

; ── Logging ─────────────────────────────────────────────────────────────────
stdout_logfile=/var/www/cuan-workflow/storage/logs/supervisor-worker.log
stderr_logfile=/var/www/cuan-workflow/storage/logs/supervisor-worker-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3

; ── Restart policy ──────────────────────────────────────────────────────────
startretries=3
startsecs=1
